# templates/cmd.service.j2
[Unit]
Description={{name}}

[Service]
Type=simple
WorkingDirectory={{path}}
ExecStart={{cmd}}

[Install]
WantedBy={{target}}
