# templates/cmd.service.j2
[Unit]
Description={{name}}

[Service]
Type=simple
WorkingDirectory={{path}}
ExecStart={{cmd}}
User={{user}}
Group={{group}}

[Install]
WantedBy={{target}}
